# network/if-up.d/static-arp
#!/bin/sh

arp 193.191.177.254 32:51:01:e9:25:5a
